# .config/systemd/user/backup-packages.service
[Unit]
Description=Backup list of installed packages

[Service]
ExecStart=/home/alex/bin/backup-packages.sh
